# deploy/agentic-core.service
[Unit]
Description=Agentic Core Engine (Rust)
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/agentic/core
EnvironmentFile=/var/www/agentic/backend/.env
Environment=RUST_LOG=info
ExecStart=/var/www/agentic/core/target/release/core
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentic-core

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/agentic/core

[Install]
WantedBy=multi-user.target
